# weather-logger/tw.psychquant.weather-logger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>tw.psychquant.weather-logger</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/che830621/bus-eta-logger/.venv/bin/python</string>
    <string>/Users/che830621/weather-logger/weather_logger.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>WEATHER_VOLUME</key><string>/Volumes/mini-2TB-SSD</string>
    <key>WEATHER_DATA_ROOT</key><string>/Volumes/mini-2TB-SSD/che-transport/weather/parquet</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/che830621/Library/Logs/weather-logger.out.log</string>
  <key>StandardErrorPath</key><string>/Users/che830621/Library/Logs/weather-logger.err.log</string>
</dict>
</plist>
